# sc4-zookeeper/resources/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh
autostart=true
autorestart=false
